# imessage-rest-ios/layout/Library/LaunchAgents/com.ericrabil.barcelonad.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.ericrabil.imessage-rest</string>
	<key>MachServices</key>
	<dict>
		<key>com.ericrabil.imessage-rest</key>
		<true/>
	</dict>
	<key>ProgramArguments</key>
	<array>
		<string>/Library/Application Support/Barcelona/imessage-rest.xpc/imessage-rest</string>
	</array>
	<key>ProcessType</key>
	<string>Standard</string>
	<key>UserName</key>
	<string>mobile</string>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
</dict>
</plist>
